# connman-reconnect.timer
[Unit]
Description=Check network

[Timer]
Persistent=true
OnUnitActiveSec=15s
OnBootSec=60s

[Install]
WantedBy=timers.target
